# scripts/crontab
# For don't receive errors, let the cron with the final whitespace
0 9 * * * sh /home/ubuntu/applications/space-flights-batch/execute-batch.sh >> /var/log/cron.log 2>&1
